# data/fastgame.service
[Unit]
Description=FastGame

[Service]
ExecStart=/usr/bin/fastgame_server

[Install]
WantedBy=multi-user.target
